# tmux-bot-cleaner.service
[Unit]
Description=Tmux Bot Cleaner PostgreSQL Monitor
After=network-online.target postgresql.service
Wants=network-online.target postgresql.service
Requires=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/lab_code
# Runs tmux_bot_cleaner.py (repo name); copy/symlink to tmux_bot_cleaner_postgrey.py on deploy if you prefer that name
ExecStart=/bin/bash -c 'cd /home/ubuntu/lab_code && exec /usr/bin/python3 tmux_bot_cleaner_postgrey.py 2>/dev/null || exec /usr/bin/python3 tmux_bot_cleaner.py'
Restart=always
RestartSec=15
StartLimitBurst=5
StartLimitIntervalSec=300
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tmux-bot-cleaner

# Environment variables
Environment=PYTHONPATH=/home/ubuntu/lab_code
Environment=PATH=/usr/local/bin:/usr/bin:/bin
Environment=HOME=/home/ubuntu

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/ubuntu/lab_code/TmuxCleaner /tmp

# Resource limits
MemoryMax=512M
CPUQuota=50%

[Install]
WantedBy=multi-user.target
